# scripts/mystocks.cron
#! /bin/bash

/var/www/mystocks/mystocks backup
